# frontend/node-frontend.supervisord.conf
[program:node-frontend]
command=TU SKOŃCZYLISMY W PIĄTEK BARDZO ZMĘCZENI
autorestart=true
priority=300
directory=/frontend
