# frontend/node-frontend.supervisord.conf
[program:node-frontend]
command=npm run dev -- --host
autorestart=true
priority=300
directory=/frontend
