# backend/unicorn-backend.supervisord.conf
[program:uvicorn-backend]
command=/run-uvicorn-backend.sh
autorestart=true
priority=200
directory=/backend