# portroute/portroute.supervisord.conf
[program:portroute]
command=/usr/local/bin/forward-ports.sh
autorestart=true
priority=100