# comfyui/comfyui.supervisord.conf
[program:comfyui]
command=python3 main.py --listen 0.0.0.0
autorestart=true
priority=100
directory=/comfyui
